# templates/cassandra.service.erb
# Managed by Puppet.  Changes are likely to be overwritten.

[Unit]
Description=Cassandra
After=network.target

[Service]
PIDFile=/var/run/<%= @service_name %>/<%= @service_name %>.pid
ExecStart=/etc/init.d/<%= @service_name %> start
ExecStop=/usr/bin/kill $MAINPID
ExecStopPost=/usr/bin/rm -f /var/lock/subsys/<%= @service_name %>
StandardOutput=journal
StandardError=journal
LimitNOFILE=100000
LimitMEMLOCK=infinity
LimitNPROC=32768
LimitAS=infinity
SuccessExitStatus=143
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
