# templates/datastax-agent.service.erb
# Managed by Puppet.  Changes are likely to be overwritten.

[Unit]
Description=DataStax Agent
After=network.target

[Service]
PIDFile=/var/run/<%= @service_name %>/<%= @service_name %>.pid
ExecStart=/etc/init.d/<%= @service_name %> start
StandardOutput=journal
StandardError=journal
LimitNOFILE=100000
LimitMEMLOCK=infinity
LimitNPROC=32768
LimitAS=infinity

[Install]
WantedBy=multi-user.target
