# templates/opscenterd.service.erb
# Managed by Puppet.  Changes are likely to be overwritten.

[Unit]
Description=Cassandra cluster manager opscenterd
After=network.target

[Service]
PIDFile=/var/run/<%= @service_name %>/<%= @service_name %>.pid
ExecStart=/etc/init.d/<%= @service_name %> start
ExecStop=/etc/init.d/<%= @service_name %> stop
StandardOutput=journal
StandardError=journal
LimitNOFILE=100000
LimitMEMLOCK=infinity
LimitNPROC=32768
LimitAS=infinity
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
