# templates/opscenterd.service.erb
# Managed by Puppet.  Changes are likely to be overwritten.

[Unit]
Description=Cassandra cluster manager opscenterd
After=network.target

[Service]
Environment=SYSTEMCTL_SKIP_REDIRECT=1
ExecStart=/etc/init.d/opscenterd start
ExecStop=/etc/init.d/opscenterd stop
PIDFile=/var/run/opscenter/opscenterd.pid
Type=forking

[Install]
WantedBy=multi-user.target
